# openrc/ssh-trap
#!/sbin/openrc-run
command=/usr/local/sbin/ssh-trap
#command_args=""
command_background=true
pidfile=/var/run/ssh-trap.pid
name="ssh-trap"

